# s6/services/localkit-homeassistant/run
#!/command/with-contenv bash
echo "Localkit: Go2rtc";
exec 2>&1
exec /usr/local/bin/go2rtc -c /var/www/html/storage/app/private/go2rtc.yml

